# overlay.d/99custom/etc/NetworkManager/dispatcher.d/50-tailscale
#!/bin/sh

ethtool -K ${DEVICE_IFACE} rx-udp-gro-forwarding on rx-gro-list off
